# .setup/build/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:laravel-worker]
command=php /var/www/localhost/htdocs/artisan queue:work --sleep=3 --tries=3
directory=/var/www/localhost/htdocs
stdout_logfile=/var/www/localhost/htdocs/storage/logs/worker.log
redirect_stderr=true
